# dockers/explorer/explorer.conf
[supervisord]
nodaemon=true

[program:explorer]
command=node /usr/src/explorer/app.js > /dev/null 2>&1 &
autostart=true
