# dockers/explorer/grabber.conf
[supervisord]
nodaemon=true

[program:grabber]
command=node /usr/src/explorer/tools/grabber.js > /dev/null 2>&1 &
autostart=true
